# deploy/visitor-api.service.example
# /etc/systemd/system/visitor-api.service
# sudo systemctl daemon-reload && sudo systemctl enable --now visitor-api

[Unit]
Description=Visitor Access Control API (FastAPI / uvicorn)
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/visitor-api/backend_python
Environment=PATH=/var/www/visitor-api/venv/bin
ExecStart=/var/www/visitor-api/venv/bin/uvicorn main:app --host 127.0.0.1 --port 8000
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
